# node_builder/modules/scheduled_reboot/scheduled-reboot.service
[Unit]
Description=Scheduled Reboot

[Service]
Type=simple
ExecStart=/usr/local/bin/graceful_reboot
